# supervisord.conf
[program:worker1]
command=python /app/worker1.py

[program:worker2]
command=python /app/worker2.py
